# deploy/ops/crm-backup.service
[Unit]
Description=CRM encrypted PostgreSQL backup and off-server copy
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
EnvironmentFile=-/opt/crm-system/secrets/backup.env
WorkingDirectory=/opt/crm-system/current
ExecStart=/usr/bin/env bash -lc '/opt/crm-system/current/deploy/backup/backup.sh | tail -n 1 | xargs /opt/crm-system/current/deploy/backup/offsite-copy.sh'
User=root
Group=root
UMask=0077

[Install]
WantedBy=multi-user.target
